# libs/nlu/doly-nlpjs/doly-nlu.service
[Unit]
Description=Doly NLU Service (NLP.js)
Documentation=https://docs.example.com/nlu
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/dolydev/libs/nlu/doly-nlpjs
Environment="NODE_ENV=production"

# 启动命令
ExecStart=/usr/bin/node /home/pi/dolydev/libs/nlu/doly-nlpjs/server.js

# 重启策略
Restart=always
RestartSec=10

# 进程管理
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=10

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=doly-nlu-service

# 资源限制
LimitNOFILE=65535
MemoryLimit=512M

[Install]
WantedBy=multi-user.target
