# libs/widgets/scripts/widget-service.service
[Unit]
Description=Widget Service for Doly LCD Display
Documentation=file:///home/pi/dolydev/libs/widgets/docs/00_设计文档.md
After=network.target syslog.target
Wants=network.target
PartOf=doly.target

[Service]
Type=simple
WorkingDirectory=/home/pi/dolydev/libs/widgets
ExecStart=/home/pi/dolydev/libs/widgets/build/widget_service --config /home/pi/dolydev/libs/widgets/config/widgets.default.json
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=widget-service
User=pi
Group=pi
LimitNOFILE=4096
LimitNPROC=4096

# 安全性设置
# 注意: 禁用 PrivateTmp 以允许与其他服务共享 /tmp 中的 ZMQ IPC socket
PrivateTmp=no
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
